# common/systemd/cppm.path
[Unit]
Description=Watch for changes in CPPM

[Path]
PathChanged=/opt/cppm/mode

[Install]
WantedBy=multi-user.target
